# .config/systemd/user/update-brightness.service
[Unit]
Description=Update brightness using the script based on wl-gammarelay-rs

[Service]
ExecStart=%h/.config/waybar/scripts/update-brightness.sh

