# .config/systemd/user/update-temperature.timer
[Unit]
Description=Run temperature update every minute

[Timer]
OnCalendar=
OnCalendar=*-*-* *:*:00
OnActiveSec=1min
OnUnitActiveSec=1min
Persistent=true

[Install]
WantedBy=timers.target

