# .os/programs/security/clamav/services/clamav-daily-scan.timer
[Unit]
Description=Run clamav daily scan

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
