# .os/programs/security/rkhunter/services/rkhunter-scan.timer
[Unit]
Description=Run rkhunter scan daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

